# systemd/ccai-jetbot.service
[Unit]
Description=CCAI JetBot ROS Patrol (Docker)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=roboat
WorkingDirectory=/home/roboat/work/ros2_ws/ccai_jetbot_ros
EnvironmentFile=-/home/roboat/work/ros2_ws/ccai_jetbot_ros/.env
Environment=CCAI_SAFE_START=0
ExecStart=/home/roboat/work/ros2_ws/ccai_jetbot_ros/scripts/host_docker_run.sh
ExecStop=/usr/bin/docker stop ccai-jetbot
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
