# systemd/ccai-jetbot.service
[Unit]
Description=CCAI JetBot ROS Patrol
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jetson
WorkingDirectory=/home/jetson/ccai_jetbot_ros
EnvironmentFile=-/home/jetson/ccai_jetbot_ros/.env
ExecStart=/home/jetson/ccai_jetbot_ros/scripts/run_patrol.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
